# supervisord.conf
[supervisord]
nodaemon=true
pidfile=/home/xray/supervisord.pid

[program:crond]
command=/usr/sbin/crond -f >> ~/cronlog
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0

